# etc/sv/vsftpd/run
#!/bin/sh
# vsftpd - Very Secure FTP Daemon
#
exec 2>&1

NAME=${PWD##*/}
DAEMON=$(command -v $NAME)
[ -x "$DAEMON" ] || exit

exec "$DAEMON" -olisten=YES
