# etc/sv/dnsmasq/finish
#!/bin/sh
[ -x /sbin/resolvconf ] && /sbin/resolvconf -d lo.dnsmasq
