# etc/sv/rpcbind/finish
#!/bin/sh
# Stop NFS client and server support services
# (jan@wespe.dev)

sv down /etc/service/rpc.* >/dev/null
