# etc/sv/rpcbind/finish
#!/bin/sh
# Stop NFS client and server support services
# (cloux@rote.ch)

sv down /etc/service/rpc.* >/dev/null
