# etc/sv/rpc.nfsd/finish
#!/bin/sh
# kill all running nfsd instances
killall -9 --quiet --exact nfsd
