# etc/sv/socklog/run
#!/bin/sh
#
# socklog - system and kernel logging services
#
exec 2>&1

NAME=${PWD##*/}
DAEMON="$(command -v $NAME)"
[ -x "$DAEMON" ] || exit

exec chpst -Unobody "$DAEMON"
