# etc/sv/hiawatha/finish
#!/bin/sh
pkill -f "$(command -v hiawatha)"
exit 0
